# systemfiles/ourapost.timer
[Unit]
Description=Daily oura post to influxdb

[Timer]
OnCalendar=*-*-* 21:00:00
Persistent=true

[Install]
WantedBy=timers.target
